# packages/agent-dispatch/launchd/org.knowledge-medium.agent-dispatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template — substitute __NODE__ (absolute node path) and __REPO__
     (absolute repo root) before installing; see the package README. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.knowledge-medium.agent-dispatch</string>

  <key>ProgramArguments</key>
  <array>
    <string>__NODE__</string>
    <string>__REPO__/packages/agent-dispatch/dist/daemon.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO__</string>

  <key>RunAtLoad</key>
  <true/>

  <!-- SuccessfulExit=false: relaunch only when the process exits
       NON-zero (a crash). The daemon exits 0 on config errors for
       exactly this reason (restarting won't fix a bad config), and
       never exits at all while merely waiting for the bridge/tab. -->
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <!-- On stop/reboot the daemon drains in-flight runs before exiting; a
       run can legitimately take minutes (watcher timeoutMs, capped at
       25m in config.ts). Cover the full cap plus grace so launchd never
       SIGKILLs a live drain (default is ~20s, which would orphan the
       child agent process and lose its reply). This costs nothing when idle —
       launchd only waits while the process is actually still running;
       system shutdown applies its own shorter global limit, where long
       runs degrade to the requeue-after-stale-sweep path. -->
  <key>ExitTimeOut</key>
  <integer>1560</integer>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/km-agent-dispatch.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/km-agent-dispatch.log</string>

  <!-- Keep PATH minimal but sufficient for local agent CLIs (claude is
       installed via its own launcher in ~/.local/bin by default) and node. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
