# packages/claude-tasks/launchd/org.knowledge-medium.claude-tasks.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template — substitute __NODE__ (absolute node path) and __REPO__
     (absolute repo root) before installing; see the package README. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.knowledge-medium.claude-tasks</string>

  <key>ProgramArguments</key>
  <array>
    <string>__NODE__</string>
    <string>__REPO__/packages/claude-tasks/dist/daemon.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO__</string>

  <key>RunAtLoad</key>
  <true/>

  <!-- SuccessfulExit=false: relaunch only when the process exits
       NON-zero (a crash). The daemon exits 0 on config errors for
       exactly this reason (restarting won't fix a bad config), and
       never exits at all while merely waiting for the bridge/tab. -->
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <!-- On stop/reboot the daemon drains in-flight runs before exiting; a
       run can legitimately take minutes (watcher timeoutMs, up to 25m).
       Give it room before launchd SIGKILLs (default is ~20s, which would
       orphan the claude child). Graceful drain still only covers runs
       that finish within this window; longer ones degrade to the
       requeue-after-stale-sweep path. -->
  <key>ExitTimeOut</key>
  <integer>120</integer>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/km-claude-daemon.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/km-claude-daemon.log</string>

  <!-- Keep PATH minimal but sufficient for `claude` (installed via its
       own launcher in ~/.local/bin by default) and node. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
